# pk4-generate-index.timer
[Unit]
Description=pk4-generate-index (deferred)
Documentation=man:pk4-generate-index

[Timer]
OnActiveSec=1m